# workflows/bootc-debs-builder/dc-firewall-seed/usr/lib/systemd/system/dc-firewall-seed.service
[Unit]
Description=Seed DaemonCores VE firewall config into pmxcfs on first boot
Documentation=https://github.com/DaemonCores/DaemonCores-VE
# Needs the pmxcfs cluster filesystem (/etc/pve) mounted by pve-cluster. It waits
# for pmxcfs to finish initialising and seeds best-effort, so it does not order
# before pve-firewall — that daemon re-reads the rules from pmxcfs when they land.
After=pve-cluster.service
Wants=pve-cluster.service
ConditionPathExists=/usr/share/daemoncores-ve/firewall/cluster.fw.tpl

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/dc-firewall-seed

[Install]
WantedBy=multi-user.target
